# dbus.socket
[Unit]
Description=D-Bus Socket Message Bus Socket

[Socket]
ListenStream=%t/dbus/user_bus_socket

